# definitions/comfyui/comfyui.container
# comfyui.container
# Quadlet container file for ComfyUI — GPU-accelerated image generation
#
# This definition targets AMD GPUs via ROCm. For NVIDIA GPUs,
# see the NVIDIA section in README.md.

[Unit]
Description=ComfyUI — GPU-Accelerated Image Generation
Documentation=https://github.com/comfyanonymous/ComfyUI
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/yanwk/comfyui-boot:rocm7
ContainerName=comfyui

# --- GPU: AMD ROCm ---
# /dev/kfd: kernel fusion driver (HSA compute)
# /dev/dri: direct rendering (GPU enumeration)
AddDevice=/dev/kfd
AddDevice=/dev/dri

# GPU group access (video=39, render=105 on Fedora/RHEL)
# Verify with: getent group video render
GroupAdd=39
GroupAdd=105

# Security settings required by ROCm runtime
SecurityLabelDisable=true
PodmanArgs=--ipc=host
PodmanArgs=--cap-add=SYS_PTRACE
PodmanArgs=--security-opt=seccomp=unconfined

# --- GPU: NVIDIA (uncomment and remove AMD section above) ---
# AddDevice=nvidia.com/gpu=all
# (no extra security settings needed for NVIDIA)

# ROCm environment — adjust HSA_OVERRIDE_GFX_VERSION for your GPU
# See README.md for common GFX version values
EnvironmentFile=./comfyui.env

# ComfyUI launch arguments
# --listen 0.0.0.0: bind to all interfaces (required for port publishing)
# --lowvram: offload model layers to CPU when VRAM is tight
# --reserve-vram 1: keep 1 GB free for system overhead
Environment="CLI_ARGS=--listen 0.0.0.0"

# Web UI port
PublishPort=8188:8188

# Shared model storage — see docs/model-storage.md
# Uses :z (shared) suffix because multiple containers may access these paths
Volume=/srv/models/comfyui:/root/ComfyUI/models:z
Volume=/srv/models/huggingface:/root/.cache/huggingface:z
Volume=/srv/models/torch:/root/.cache/torch:z

# User directory — persists ComfyUI-Manager cache, workflows, assets DB, and log.
# Without this, the Manager re-fetches ComfyRegistry data (133 pages) and
# 5 GitHub JSON files on every startup.
Volume=/srv/containers/comfyui/user:/root/ComfyUI/user:Z

# Per-application data — inputs, outputs
Volume=/srv/containers/comfyui/input:/root/ComfyUI/input:Z
Volume=/srv/containers/comfyui/output:/root/ComfyUI/output:Z

# Custom nodes — installed extensions persist here
Volume=/srv/containers/comfyui/custom_nodes:/root/ComfyUI/custom_nodes:Z

[Service]
Restart=on-failure
# First run pulls the image and installs ComfyUI + dependencies
TimeoutStartSec=900

[Install]
WantedBy=default.target
